# deployment/twitter-analyzer.service
[Unit]
Description=Twitter Profile Analyzer API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/followfeed
Environment=PYTHONPATH=/root/followfeed
Environment=PATH=/usr/bin:/usr/sbin:/bin:/sbin:/root/followfeed/xfeed/bin:$PATH

# Load environment variables from .env file
EnvironmentFile=/root/followfeed/.env

# Use API_PORT environment variable if set, otherwise default to 8000
ExecStart=/bin/bash -c '/root/followfeed/xfeed/bin/uvicorn api.app:app --host 0.0.0.0 --port 8000 --workers 4 --timeout-keep-alive 120 --log-level warning --limit-concurrency 100 --backlog 200 --proxy-headers'

# Restart configuration
Restart=always
RestartSec=10
StartLimitIntervalSec=0
StartLimitBurst=0

# Resource limits
CPUQuota=50%
MemoryMax=4G

# Logging
StandardOutput=append:/var/log/twitter-analyzer/api.log
StandardError=append:/var/log/twitter-analyzer/api-error.log

[Install]
WantedBy=multi-user.target 